# root/etc/services.d/track/run
#!/usr/bin/with-contenv sh

cd /track

s6-setuidgid guest /bin/sh -c "exec node wsserver.js;"
